# templates/default/onlinefs.service.erb
[Unit]
Description = Online feature store service.
After = syslog.target network.target remote-fs.target docker.service <%= @deps %>
Requires = docker.service

[Service]
ExecStart = /usr/bin/docker run --rm --name %n \
           -v <%= node['onlinefs']['home'] %>:/onlinefs \
           -v <%= node['onlinefs']['data_volume']['root_dir']%>:<%= node['onlinefs']['data_volume']['root_dir']%> \
           -v <%= @crypto_dir %>:/certs \
           -p <%= node['onlinefs']['monitoring'] %>:<%= node['onlinefs']['monitoring'] %> \
           -e ONLINEFS_HOME=/onlinefs \
           docker.hops.works/onlinefs:<%= node['onlinefs']['version'] %>
ExecStop = /usr/bin/docker stop %n 

Restart = on-failure
RestartSec = 2

[Install]
WantedBy = multi-user.target
